# conf/systemd.path
[Unit]
Description="Lauch 'faircamp --build-dir' on content change"

[Path]
PathModified=__INSTALL_DIR__/content
Unit=__APP__.service

[Install]
WantedBy=multi-user.target